# templates/default/systemd.service.erb
[Unit]
Description="<%= @name %>"

[Service]
Type=simple
ExecStart=/bin/java -Djava.library.path=<%= @path %>/DynamoDBLocal_lib -jar <%= @path %>/DynamoDBLocal.jar -port <%= @port %> -dbPath /data >> <%= @log_dir %>/<%= @port %>.log 2>&1
User=<%= @user %>
Group=<%= @group %>
Type=simple

[Install]
WantedBy=multi-user.target

